# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/workspace/logs/supervisord.log
pidfile=/tmp/supervisord.pid
user=root
environment=OPENCLAW_HOME="/workspace/.openclaw",OPENCLAW_GATEWAY_MODE="standalone"

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
startretries=3
priority=1
stdout_logfile=/workspace/logs/sshd.log
stderr_logfile=/workspace/logs/sshd.err

[program:llm]
command=/scripts/start_llm.sh
directory=/workspace
autostart=true
autorestart=unexpected
startretries=3
priority=2
stdout_logfile=/workspace/logs/llm.log
stderr_logfile=/workspace/logs/llm.err

[program:openclaw]
command=/scripts/start_openclaw.sh
directory=/workspace
autostart=true
autorestart=unexpected
startretries=3
priority=3
stopasgroup=true
killasgroup=true
stopwaitsecs=10
stdout_logfile=/workspace/logs/openclaw.log
stderr_logfile=/workspace/logs/openclaw.err

[program:jupyter]
command=/scripts/start_jupyter.sh
directory=/workspace
autostart=true
autorestart=unexpected
startretries=3
priority=4
stdout_logfile=/workspace/logs/jupyter.log
stderr_logfile=/workspace/logs/jupyter.err
