# supervisord.conf
[supervisord]
nodaemon=true
logfile=/workspace/logs/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
startretries=3
priority=1
stdout_logfile=/workspace/logs/sshd.log
stderr_logfile=/workspace/logs/sshd.err

[program:jupyter]
command=/opt/venv/bin/jupyter lab --ip=0.0.0.0 --port=8888 --no-browser --allow-root --NotebookApp.token=''
directory=/workspace
autostart=true
autorestart=unexpected
startretries=3
priority=2
stdout_logfile=/workspace/logs/jupyter.log
stderr_logfile=/workspace/logs/jupyter.err

[program:llm]
command=/scripts/start_llm.sh
directory=/workspace
autostart=true
autorestart=unexpected
startretries=3
priority=3
stdout_logfile=/workspace/logs/llm.log
stderr_logfile=/workspace/logs/llm.err

[program:openclaw]
command=/bin/bash -lc 'openclaw gateway run --port 18789 --allow-unconfigured'
directory=/workspace
autostart=true
autorestart=unexpected
startretries=3
priority=4
stdout_logfile=/workspace/logs/openclaw.log
stderr_logfile=/workspace/logs/openclaw.err
